# salt/app/files/supervisor.conf
{% set IP = salt['network.interface_ip']('eth0')  %}
[program:app]
command=/home/app/env/bin/python /home/app/env/bin/gunicorn  app:app -b {{ IP }}:8000
user=app
directory=/home/app/code
