# salt/app/files/supervisor.conf
[program:app]
command=/home/app/env/bin/python /home/app/env/bin/gunicorn  app:app -b 0.0.0.0:8000
user=app
directory=/home/app/code
